# pg_metrics.timer
[Unit]
Description=Executa pg_metrics a cada 30 segundos

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=1s

[Install]
WantedBy=timers.target
